# supervisord-go-agent.conf
[program:go-agent]
command=/usr/share/go-agent/agent.sh
autorestart=true

